# scripts/com.dotatrader.dailycycle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dotatrader.dailycycle</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-l</string>
        <string>/Users/yanpushkaryov/Downloads/dota_trader_v2/scripts/daily_paper_cycle.sh</string>
    </array>

    <!-- Continuous-режим: запуск каждые 10 минут (600 сек).
         Тяжёлые PandaScore-вызовы троттлятся внутри скриптов, так что
         частый запуск безопасен — лишние прогоны почти мгновенны. -->
    <key>StartInterval</key>
    <integer>600</integer>

    <!-- Сразу выполнить один раз при загрузке launchd-агента -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Если Mac был выключен/спал в момент запуска — выполнит при следующей загрузке -->
    <key>StandardOutPath</key>
    <string>/Users/yanpushkaryov/Downloads/dota_trader_v2/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yanpushkaryov/Downloads/dota_trader_v2/logs/launchd_stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/yanpushkaryov/Downloads/dota_trader_v2</string>
</dict>
</plist>
